# SPM_AgentService_Linux/SPM_AgentService_Linux/bin/Builded/spm-agent/etc/systemd/system/spm-agent.service
[Unit]
Description=SPM Monitoring Agent Service
[Service]
ExecStart=/bin/SPM_AgentService_Linux
[Install]
WantedBy=multi-user.target
